# build/common/root/supervisor-remotion.conf
[program:remotion]
command=sh -c "echo 'No render script configured - mount your remotion project to /config or /data'; sleep infinity"
stdout_logfile=/config/supervisor-remotion.log
stderr_logfile=/config/supervisor-remotion.err.log
autostart=true
autorestart=true
priority=5
